# skills/cs-weekly/cron/com.avada.cs-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.avada.cs-weekly</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/avada/CSL/skills/cs-weekly/cron/run-weekly.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/avada/CSL</string>

  <!-- Every Monday at 09:00 local time. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <!-- If the Mac was asleep/off at 09:00 Mon, launchd runs it as soon as it wakes. -->
  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/cs-weekly.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/cs-weekly.log</string>
</dict>
</plist>
