# skills/dfy-tracker/cron/com.avada.dfy-tracker-monthly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.avada.dfy-tracker-monthly</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/avada/CSL/skills/dfy-tracker/cron/run-monthly.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/avada/CSL</string>

  <!-- 2nd of every month at 15:00 (3PM) local time. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Day</key>
    <integer>2</integer>
    <key>Hour</key>
    <integer>15</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <!-- Machine off/asleep at 15:00 on the 2nd → launchd skips this run.
       Run manually the next day: bash run-monthly.sh -->
  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/dfy-tracker-monthly.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/dfy-tracker-monthly.log</string>
</dict>
</plist>
