# skills/mine-chat-faqs/cron/com.avada.mine-faqs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.avada.mine-faqs</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/avada/CSL/skills/mine-chat-faqs/cron/run-weekly.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/avada/CSL</string>

  <!-- Every Monday at 16:00 local time. Mines the previous full Mon→Sun week. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>16</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <!-- launchd runs a missed calendar job once on next wake; it does not wake the Mac. -->
  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/mine-faqs-weekly.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/mine-faqs-weekly.log</string>
</dict>
</plist>
