# skills/mine-chat-faqs/cron/com.avada.mine-faqs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.avada.mine-faqs</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/avada/CSL/skills/mine-chat-faqs/cron/run-weekly.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/avada/claw-weebhook-crisp-chat</string>

  <!-- Every Tuesday at 10:00 local time. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>2</integer>
    <key>Hour</key>
    <integer>10</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <!-- If the Mac was asleep/off at 10:00 Tue, run as soon as it wakes. -->
  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/mine-faqs-weekly.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/mine-faqs-weekly.log</string>
</dict>
</plist>
